# packaging/systemd/blizzard-runner.service
# The blizzard-runner daemon — the supervisor: the stateless REAP -> PULL -> FILL
# -> ADVANCE reconciliation loop behind a local API (glossary "supervisor",
# design/cli.md, D-061). Installed as the runner half of the colocated
# single-machine pair (product/mvp.md journey); the hub half is
# blizzard-hub.service. Setup + the reboot-recovery contract: docs/deployment.md.
[Unit]
Description=Blizzard runner — the supervisor loop + local API (blizzard-runner, D-061)
Documentation=https://github.com/paul-gross/blizzard/blob/master/docs/deployment.md
# Colocated ordering: prefer the hub up first, but the runner reaches the hub
# outbound-only (D-012) and keeps working while the hub is unreachable
# (design/cli.md), so this is a convenience, not a hard dependency.
After=network-online.target blizzard-hub.service
Wants=network-online.target

[Service]
Type=exec
# Offline schema reconcile before the daemon opens the store (D-099); the daemon
# refuses to start on a revision mismatch. Self-healing across a wheel upgrade.
ExecStartPre=/opt/blizzard/venv/bin/blizzard-runner migrate --dir /var/lib/blizzard/runner
ExecStart=/opt/blizzard/venv/bin/blizzard-runner host --dir /var/lib/blizzard/runner
# Optional overrides (workspace root, hub URL, harness binary) — the durable ones
# live in the runner's own config.toml, seeded by `blizzard runner init`. No forge
# or PM credentials belong here (D-084): a runner never holds them.
EnvironmentFile=-/etc/blizzard/runner.env
User=blizzard
Group=blizzard
# Shares /var/lib/blizzard with the hub unit; runner state lives under runner/.
StateDirectory=blizzard
WorkingDirectory=/var/lib/blizzard/runner
# The reboot/crash-recovery mechanism (docs/deployment.md): after a kill -9 or a
# machine reboot, systemd brings the supervisor back and its FIRST tick is REAP —
# it reaps the leases stranded by the crash, re-reads its environment bindings from
# the store, and continues each chunk at the node the hub last recorded; nothing is
# worked twice (proven by tests/crash/test_kill9_sweep.py::test_kill9_runner_daemon_mid_flight).
#
# A GRACEFUL restart (`systemctl restart`, or stop→start on a wheel upgrade) is stronger:
# the SIGTERM lets the daemon run its shutdown path, which marks every in-flight lease
# with a durable resume-intent (D-082). Its first startup tick then RESUMEs each marked
# session in place — same lease/epoch/session, only the pid rewritten, no retry consumed —
# so adopting a newly-merged runner build no longer discards in-flight agent context. The
# default KillMode (control-group) SIGTERMs the workers alongside the daemon; RESUME kills
# the marked worker first so it never runs two processes on one session (D-049). A `kill -9`
# during the re-attach itself still degrades to the ordinary crash contract above (REAP
# requeues) — the recovery boundaries are the crash-point sweep's, not a stronger claim than
# it arms; see design/runner/loop.md and steps.py `_resume_in_place`.
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
